# systemd/pifan.service
[Unit]
Description=Raspberry Pi fan control

[Service]
Type=simple
User=admin
ExecStart=/usr/sbin/pifan -start 66 -stop 60 -timeout 30 -thermal /sys/class/thermal/thermal_zone0/temp -gpio 2
Restart=always
KillSignal=SIGQUIT

[Install]
WantedBy=default.target
